# sonarr.ini
[program:rtorrent]
user=
command=/usr/bin/mono /opt/sonarr/NzbDrone.exe -nobrowser -data /config/
