# sonarr.ini
[program:sonarr]
user=
Environment=MONO_TLS_LEGACY
command=/usr/bin/mono /opt/sonarr/NzbDrone.exe /nobrowser /data=/config/
